# config/ollama.service.example
[Unit]
Description=Ollama AI Server with Open WebUI via Docker Compose
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target
Conflicts=localai.service

[Service]
Type=oneshot
WorkingDirectory=/opt/ollama
Environment=COMPOSE_PROJECT_NAME=ollama
ExecStartPre=/usr/bin/docker compose pull --quiet
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
